# celebdating/templates/default/celeb_backend_service.conf.erb
description "Celebrity Dating"

start on runlevel [234]
stop on runlevel [0156]

kill timeout 1800
limit nofile 65536 65536

setuid <%= @node[:celebdating][:runuser] %>
setgid neon

chdir <%= @repo_root %>
env VIRTUAL_ENV=<%= @repo_root %>/.pyenv
env PATH="${VIRTUAL_ENV}/bin:$PATH"
exec ${VIRTUAL_ENV}/bin/python <%= @repo_root %>/server/api/controllers.py --db_host <%= @db[:host]%> --db_user <%= @db[:username] %> --db_port <%= db[:port] %> --db_password <%= @db[:password] %> --face_model <%= face_model_file %> --celebreity_model <%= celebrity_model_file %> >> <%= @node[:celebdating][:log_dir] %>/celeb_backend.stdout 2> <%= @node[:celebdating][:log_dir]%>/celeb_backend.stderr

respawn